# deploy/course-assistant.service
[Unit]
Description=课程助手 Web 服务
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/lecture_tools
ExecStart=/root/lecture_tools/bin/server
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
